# supervisor.conf
[program:gunicorn]
command=/home/capua/venv/bin/gunicorn app:app --config gunicorn.conf.py
directory=/home/capua/RemessaDoc
autostart=true
autorestart=true
user=capua